# utils/llamacpp.service
[Unit]
Description=llama.cpp inference server
Documentation=https://github.com/ggml-org/llama.cpp
After=network-online.target nvidia-persistenced.service
Wants=network-online.target nvidia-persistenced.service

[Service]
Type=simple
User=llama
Environment=CUDA_VISIBLE_DEVICES=SUB_CUDA_DEVICE_HERE
Group=llama
WorkingDirectory=SUB_LLAMA_PATH_HERE

ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do nvidia-smi -L > /dev/null 2>&1 && exit 0; sleep 1; done; exit 1'
ExecStart=SUB_LLAMA_PATH_HERE/build/bin/llama-server \
    -m SUB_MODEL_DIR_HERE/SUB_MODEL_HERE \
    --n-gpu-layers SUB_GPU_LAYERS_HERE \
    --ctx-size SUB_CTX_SIZE_HERE \
    --parallel SUB_SLOTS_HERE \
    --flash-attn on \
    --jinja \
    --host 0.0.0.0 \
    --port 8502 \
    --api-key SUB_API_KEY_HERE \
    --metrics \
    --log-timestamps \
    --log-prefix \
    --cache-ram SUB_PROMPT_CACHE_SIZE_HERE

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=SUB_LLAMA_PATH_HERE SUB_MODEL_DIR_HERE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-server

[Install]
WantedBy=multi-user.target
